# files/systemd/cosmic-generic-startup-hook.path
[Unit]
Wants=cosmic-generic-startup-hook.service

[Path]
PathChanged=/tmp/cmdline

[Install]
WantedBy=multi-user.target
