# files/systemd/cosmic-generic-startup-hook.path
[Unit]
Wants=cosmic-generic-startup-hook.service

[Path]
PathChanged=/tmp/cmdline
